# launch_agents/com.hiastro.business-dashboard-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hiastro.business-dashboard-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin; WORK_ROOT=/Users/yashs/Documents/WorkDirectory; DASHBOARD_ROOT=$WORK_ROOT/hiastro-business-dashboard; PYTHON=$WORK_ROOT/.venv/bin/python; GIT=/usr/bin/git; cd "$WORK_ROOT" || exit 1; for attempt in 1 2 3; do "$PYTHON" "$DASHBOARD_ROOT/scripts/build_dashboard_data.py" &amp;&amp; break; if [ "$attempt" = "3" ]; then exit 1; fi; echo "Dashboard refresh failed on attempt $attempt; retrying in 90 seconds..."; sleep 90; done; cd "$DASHBOARD_ROOT" || exit 1; "$GIT" fetch origin main; if ! "$GIT" diff --quiet -- data/dashboard_data.json; then "$GIT" add data/dashboard_data.json; "$GIT" commit -m "Refresh dashboard data"; "$GIT" push origin main; else echo "No dashboard data changes to commit."; fi</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yashs/Documents/WorkDirectory/hiastro-business-dashboard</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/yashs/Documents/WorkDirectory/hiastro-business-dashboard/logs/daily-refresh.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/yashs/Documents/WorkDirectory/hiastro-business-dashboard/logs/daily-refresh.err.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
